# sen2agri-processors/fmask/sen2agri-fmask.timer
[Unit]
Description=FMask runner

[Timer]
Unit=sen2agri-fmask.service
OnCalendar=minutely
Persistent=true

[Install]
WantedBy=timers.target
