# era5-downloader/dist/sen2agri-era5-downloader.service
[Unit]
Description=Retrieves ERA5 Weather products
After=docker.service
Requires=docker.service

[Service]
Type=simple
ExecStart=/usr/share/sen2agri/era5_weather_downloader/weather_launcher.py
User=sen2agri-service
