# sen2agri-processors/fmask/dist/sen2agri-fmask.service
[Unit]
Description=Runs FMask on L2A products
After=docker.service
Requires=docker.service

[Service]
Type=simple
ExecStart=/usr/share/sen2agri/fmask/fmask_launcher.py
User=sen2agri-service
